# tools/mpris-proxy.service.in
[Unit]
Description=Bluetooth to MPRIS proxy

[Service]
Type=simple
ExecStart=@bindir@/mpris-proxy
NotifyAccess=main
LimitNPROC=1

[Install]
WantedBy=bluetooth.target
